# package/openrc/neko
#!/sbin/openrc-run

name="Neko"
description="An(other) uptime monitor"
pidfile="${pidfile:-"/run/${SVCNAME}.pid"}"
command_user="root"

supervisor=supervise-daemon
command="/usr/bin/neko"
command_args="$ARGS"
command_background="yes"
capabilities="^cap_net_raw"

start_stop_daemon_args="--stdout /var/log/${SVCNAME}/${SVCNAME}.log \
    --stderr /var/log/${SVCNAME}/${SVCNAME}.log"                 


depend() {
    need net
    after firewall
}

start_pre() {
    checkpath --directory --owner $command_user:$command_user --mode 0755 \
        /var/log/$RC_SVCNAME
}